# scripts/runit/uwsgi/run
#!/bin/bash
exec uwsgi --emperor /var/www/app/scripts/runit/uwsgi/